# scripts/etc/uci-defaults/99-permissions.sh
#!/bin/sh
chmod 700 /etc/dropbear
echo 'chmod 700 /etc/dropbear'
chmod 600 /etc/dropbear/authorized_keys
echo 'chmod 600 /etc/dropbear/authorized_keys'
